# IC_team5/yocto/meta-wifi/recipes-connectivity/source/route-i-wlan0.service
[Unit]
Description=WPA Supplicant
Requires=udhcpc-i-wlan0.service
After=udhcpc-i-wlan0.service
Wants=network.target
[Service]
ExecStart=/sbin/route add -n 224.224.224.243 dev wlan0
RemainAfterExit=yes
Restart=on-failure
RestartSec=3


[Install]
WantedBy=multi-user.target
